# roles/monitoring/alertmanager/templates/alertmanager.service.j2
[Unit]
Description=Prometheus Alertmanager
Documentation=https://prometheus.io/docs/alerting/alertmanager/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ alertmanager_user }}
Group={{ alertmanager_group }}
ExecStart={{ alertmanager_install_dir }}/alertmanager \
    --config.file={{ alertmanager_config_dir }}/alertmanager.yml \
    --storage.path={{ alertmanager_data_dir }} \
    --web.listen-address={{ alertmanager_listen_address }} \
{% if alertmanager_cluster_listen_address %}
    --cluster.listen-address={{ alertmanager_cluster_listen_address }} \
{% endif %}
{% if alertmanager_external_url %}
    --web.external-url={{ alertmanager_external_url }} \
{% endif %}

ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ alertmanager_data_dir }}

[Install]
WantedBy=multi-user.target
