# roles/monitoring/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ node_exporter_user }}
Group={{ node_exporter_group }}
ExecStart={{ node_exporter_install_dir }}/node_exporter \
    --web.listen-address={{ node_exporter_listen_address }} \
    --web.telemetry-path={{ node_exporter_web_telemetry_path }} \
{% if node_exporter_textfile_dir is defined %}
    --collector.textfile.directory={{ node_exporter_textfile_dir }} \
{% endif %}
{% for collector in node_exporter_enabled_collectors %}
    --collector.{{ collector }} \
{% endfor %}
{% for collector in node_exporter_disabled_collectors %}
    --no-collector.{{ collector }} \
{% endfor %}
{% for arg in node_exporter_extra_args %}
    {{ arg }} \
{% endfor %}

Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths={{ node_exporter_textfile_dir }}

[Install]
WantedBy=multi-user.target
